# lila-engine.socket
[Unit]
Description=lila-engine socket

[Socket]
ListenStream=/run/lila-engine.socket

[Install]
WantedBy=sockets.target
